# mein-assistent-lena-imap-poller.service
[Unit]
Description=Lena IMAP Poller — pollt lena@herbertgruppe.com und erstellt Paperclip-Issues
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mein-assistent
EnvironmentFile=/opt/mein-assistent/.env
Environment="LENA_IMAP_LOG_FILE=/var/log/lena-imap-poller.log"
ExecStart=/opt/mein-assistent/venv/bin/python3 /opt/mein-assistent/lena_imap_poller.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lena-imap-poller

# Sicherheitseinstellungen
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
