# plaud-poller.service
[Unit]
Description=Plaud Poller — pollt Plaud-Account(s) und erstellt Paperclip-Issues
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=plaud
Group=plaud
WorkingDirectory=/opt/mein-assistent
EnvironmentFile=/opt/mein-assistent/.env
Environment="PLAUD_DB_PATH=/var/lib/plaud/state.db"
Environment="PLAUD_LOG_FILE=/var/log/plaud/plaud-poller.log"
ExecStart=/opt/mein-assistent/venv/bin/python3 /opt/mein-assistent/plaud_poller.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=plaud-poller

# Verzeichnisse — systemd erstellt sie mit korrekten Rechten (plaud:plaud)
StateDirectory=plaud
LogsDirectory=plaud

# systemd-Sandboxing (Least Privilege — HBE-854)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
